# scripts/jarvis-briefing.service
[Unit]
Description=ARIA morning briefing nudge (optional systemd timer)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=JARVIS_PORT=8765
ExecStart=/bin/sh -c 'curl -sf "http://127.0.0.1:${JARVIS_PORT}/api/briefing?launch=1" >/dev/null || notify-send -a Jarvis "ARIA" "Good morning — open ARIA for today'\''s briefing."'

[Install]
WantedBy=default.target
